# scripts/com.horizon.daily.icloud.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.horizon.daily.icloud</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/chenxin/code/ai-news-radar/horizon/scripts/run-and-export-icloud.sh</string>
    <string>12</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/chenxin/code/ai-news-radar/horizon</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>20</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/chenxin/code/ai-news-radar/horizon/logs/horizon-icloud.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/chenxin/code/ai-news-radar/horizon/logs/horizon-icloud.err.log</string>
</dict>
</plist>
